# deploy/sourceannotator-backup.service
[Unit]
Description=ScholarMark backup snapshot
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/app
ExecStart=/bin/bash /opt/app/deploy/backup-data.sh
